# deploy/ksaunibliss-webhook.service
[Unit]
Description=KsauniBliss Webhook Server
After=network.target

[Service]
Type=simple
User=u574849695
WorkingDirectory=/home/u574849695/ksaunibliss/deploy/webhook-server
ExecStart=/usr/bin/node webhook-server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=WEBHOOK_PORT=3001
Environment=WEBHOOK_SECRET=your-secret-key-here
Environment=DEPLOY_SCRIPT=/home/u574849695/ksaunibliss/deploy/webhook-deploy.sh
Environment=DOMAIN_NAME=ksaunibliss.com

# Logging
StandardOutput=append:/var/log/ksaunibliss-webhook.log
StandardError=append:/var/log/ksaunibliss-webhook-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/u574849695/ksaunibliss
ReadWritePaths=/var/log
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
